# tweet_scraper/launchd/com.user.tweet_scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.tweet_scraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/user/Documents/Code/SLT/tweet_scraper/scraper.py</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/user/Documents/Code/SLT/tweet_scraper/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/user/Documents/Code/SLT/tweet_scraper/output.log</string>
</dict>
</plist> 